# pkgbuild/pkgbuild.install
## arg 1:  the new package version
#pre_install() {
	# do something here
#}

## arg 1:  the new package version
post_install() {
    if [ ! -d "/run/miso" ]; then
        systemctl enable pre-systemd-swap.service
        systemctl enable systemd-swap.service
    fi
}
